# cti_parser/outbound/outbound.json
{
  "$schema": "../schemes/iocs-schema.json",
  "$id": "iocs.json",
  "reports": [
    {
      "url": "https://www.welivesecurity.com/en/eset-research/sednit-reloaded-back-trenches/",
      "created": "10/03/2026",
      "tags": [
        "Sednit",
        "APT28",
        "SlimAgent",
        "BeardShell",
        "Russia",
        "Ukraine"
      ],
      "iocs": [
        {
          "id": "ioc-1",
          "type": "file",
          "malign": true,
          "hash": "sha1",
          "value": "5603E99151F8803C13D48D83B8A64D071542F01B",
          "file_name": "eapphost.dll",
          "file_type": "dynamic library",
          "platform": "Windows",
          "capabilities": "Keylogging, screen capture, clipboard data collection",
          "detection": "Win64/Spy.KeyLogger.LS",
          "comment": "SlimAgent keylogger"
        },
        {
          "id": "ioc-2",
          "type": "file",
          "malign": true,
          "hash": "sha1",
          "value": "6D39F49AA11CE0574D581F10DB0F9BAE423CE3D5",
          "file_name": "tcpiphlpsvc.dll",
          "file_type": "dynamic library",
          "platform": "Windows",
          "capabilities": "Remote access, data collection",
          "detection": "Win64/BeardShell.A",
          "comment": "BeardShell trojan"
        },
        {
          "id": "ioc-3",
          "type": "url",
          "malign": false,
          "value": "https://icedrive.net/",
          "comment": "BeardShell C2 infrastructure, legit service abused by the malware"
        },
        {
          "id": "ioc-4",
          "type": "url",
          "malign": false,
          "value": "https://filen.com/",
          "comment": "Covenant C2 infrastructure, legit service abused by the malware"
        },
        {
          "id": "ioc-5",
          "type": "process",
          "malign": false,
          "value": "taskhost.exe",
          "comment": "BeardShell only executes in taskhost.exe or taskhostw.exe"
        },
        {
          "id": "ioc-6",
          "type": "process",
          "malign": false,
          "value": "explorer.exe",
          "comment": "SlimAgent only executes in explorer.exe"
        }
      ],
      "ttps": [
        {
          "id": "T1583.006",
          "comment": "Acquire Infrastructure: Web Services - BeardShell relies on Icedrive cloud storage. Covenant relies on Filen cloud storage."
        },
        {
          "id": "T1587.001",
          "comment": "Develop Capabilities: Malware - BeardShell and SlimAgent are custom malware."
        },
        {
          "id": "T1059.001",
          "comment": "Command and Scripting Interpreter: PowerShell - BeardShell executes PowerShell commands."
        },
        {
          "id": "T1129",
          "comment": "Shared Modules - BeardShell and SlimAgent are full-fledged DLL files."
        },
        {
          "id": "T1546.015",
          "comment": "Event Triggered Execution: Component Object Model Hijacking - BeardShell and SlimAgent are made persistent by hijacking COM objects."
        },
        {
          "id": "T1027",
          "comment": "Obfuscated Files or Information - BeardShell Icedrive token decryption is obfuscated."
        },
        {
          "id": "T1140",
          "comment": "Deobfuscate/Decode Files or Information - BeardShell decrypts its strings."
        },
        {
          "id": "T1480",
          "comment": "Execution Guardrails - BeardShell only executes in taskhost.exe or taskhostw.exe. SlimAgent only executes in explorer.exe."
        },
        {
          "id": "T1564",
          "comment": "Hide Artifacts - SlimAgent logs are written into a hidden file."
        },
        {
          "id": "T1082",
          "comment": "System Information Discovery - BeardShell sends a fingerprint of the compromised machine."
        },
        {
          "id": "T1005",
          "comment": "Data from Local System - BeardShell, Covenant, and SlimAgent collect data from a compromised machine."
        },
        {
          "id": "T1056.001",
          "comment": "Input Capture: Keylogging - SlimAgent performs keylogging."
        },
        {
          "id": "T1113",
          "comment": "Screen Capture - SlimAgent captures screenshots of the compromised machine."
        },
        {
          "id": "T1115",
          "comment": "Clipboard Data - SlimAgent collects clipboard data."
        },
        {
          "id": "T1001",
          "comment": "Data Obfuscation - BeardShell exfiltrates data in fake images."
        },
        {
          "id": "T1071.001",
          "comment": "Application Layer Protocol: Web Protocols - BeardShell and Covenant use HTTPS for C&C."
        },
        {
          "id": "T1102",
          "comment": "Web Service - BeardShell gets commands from Icedrive. Covenant gets commands from Filen."
        },
        {
          "id": "T1573.002",
          "comment": "Encrypted Channel: Asymmetric Cryptography - BeardShell communications with Icedrive are encrypted using HTTPS. Covenant communications with its controller uses RSA-encrypted session keys."
        },
        {
          "id": "T1567",
          "comment": "Exfiltration Over Web Service - BeardShell exfiltrates data to Icedrive. Covenant exfiltrates data to Filen."
        }
      ],
      "queries": [
        {
          "id": "query-01",
          "type": "kql",
          "value": "",
          "iocs": [
            "ioc-5",
            "ioc-6"
          ],
          "comment": "SlimAgent and BeardShell execution, dll injection"
        },
        {
          "id": "query-02",
          "type": "yara",
          "value": "",
          "iocs": [
            "ioc-1"
          ],
          "comment": "SlimAgent 2024 file detection by specific function names. content:\"FN_log_window_text\""
        },
        {
          "id": "query-03",
          "type": "yara",
          "value": "",
          "iocs": [
            "ioc-2"
          ],
          "comment": "Modified Covenant file detection by specific function names. content:\"GruntUUID\" and content:\"filenClient\""
        }
      ],
      "xrefs": [
        {
          "id": "xref-01",
          "type": "creates",
          "parents": [],
          "children": [],
          "comment": ""
        }
      ]
    },
    {
      "url": "https://research.checkpoint.com/2026/handala-hack-unveiling-groups-modus-operandi/",
      "created": "12/03/2026",
      "tags": [
        "Handala",
        "Void Manticore",
        "Red Sandstorm",
        "Banished Kitten",
        "Iran",
        "Israel",
        "Albania"
      ],
      "iocs": [
        {
          "id": "ioc-1",
          "type": "ip",
          "malign": true,
          "value": "107.189.19.52",
          "comment": "C2 server for additional payload retrieval"
        },
        {
          "id": "ioc-2",
          "type": "ip",
          "malign": false,
          "value": "169.150.227.0/24",
          "comment": "Commercial VPN segment used for egress to Israel targets"
        },
        {
          "id": "ioc-3",
          "type": "ip",
          "malign": false,
          "value": "146.185.219.235",
          "comment": "Israeli VPN node used for egress"
        },
        {
          "id": "ioc-4",
          "type": "process",
          "malign": false,
          "value": "NetBird",
          "comment": "Zero-trust mesh network deployed for internal lateral movement"
        },
        {
          "id": "ioc-5",
          "type": "process",
          "malign": false,
          "value": "VeraCrypt",
          "comment": "Disk encryption utility downloaded for destructive encryption"
        },
        {
          "id": "ioc-6",
          "type": "process",
          "malign": true,
          "value": "rundll32.exe comsvcs.dll MiniDump",
          "comment": "LSASS credential dumping pre-wiping"
        },
        {
          "id": "ioc-7",
          "type": "process",
          "malign": true,
          "value": "wmic.exe process call create",
          "comment": "Copies registry hives from VSS (SYSTEM hive)"
        }
      ],
      "ttps": [
        {
          "id": "T1190",
          "comment": "Exploit Public-Facing Application - Compromised VPN accounts for initial access."
        },
        {
          "id": "T1133",
          "comment": "External Remote Services - RDP for hands-on lateral movement."
        },
        {
          "id": "T1570",
          "comment": "Lateral Tool Transfer - NetBird deployment via browser download."
        },
        {
          "id": "T1003.001",
          "comment": "OS Credential Dumping: LSASS Memory - comsvcs.dll MiniDump."
        },
        {
          "id": "T1087.002",
          "comment": "Account Discovery: Domain Accounts - ADRecon PowerShell."
        },
        {
          "id": "T1485",
          "comment": "Data Destruction - Handala Wiper (MBR/file overwrite)."
        },
        {
          "id": "T1561.002",
          "comment": "Disk Wipe: Disk Structure Wipe - Custom MBR wiping."
        },
        {
          "id": "T1490",
          "comment": "Inhibit System Recovery - Multiple parallel wiping methods."
        },
        {
          "id": "T1486",
          "comment": "Data Encrypted for Impact - VeraCrypt disk encryption."
        },
        {
          "id": "T1059.001",
          "comment": "Command and Scripting Interpreter: PowerShell - AI-assisted wiping script."
        },
        {
          "id": "T1480",
          "comment": "Execution Guardrails - Domain Admin creds validation pre-impact."
        },
        {
          "id": "T1078",
          "comment": "Valid Accounts - Compromised DA credentials for GPO deployment."
        }
      ],
      "queries": [
        {
          "id": "query-1",
          "type": "kql",
          "value": "Process | where ProcessImageName in~('rundll32.exe', 'wmic.exe', 'powershell.exe') | where ProcessCommandLine contains 'comsvcs.dll' or ProcessCommandLine contains 'MiniDump' or ProcessCommandLine contains 'ADRecon'",
          "iocs": [
            "ioc-6",
            "ioc-7"
          ],
          "comment": "Pre-wiping recon: LSASS dump, registry copy, ADRecon execution"
        },
        {
          "id": "query-2",
          "type": "kql",
          "value": "DeviceNetworkEvents | where RemoteIP in ('107.189.19.52', '169.150.227.0/24', '146.185.219.235') | where ActionType == 'ConnectionSuccess'",
          "iocs": [
            "ioc-1",
            "ioc-2",
            "ioc-3"
          ],
          "comment": "Suspicious C2/VPN connections from Iranian actors"
        }
      ],
      "xrefs": [
        {
          "id": "xref-1",
          "type": "detects",
          "parents": [
            "query-1",
            "query-2"
          ],
          "children": [
            "ioc-1"
          ],
          "comment": "Compromised VPN credentials -> RDP lateral -> NetBird deployment -> Recon (ADRecon, LSASS dump) -> Could wait for months -> Disable Defender -> GPO wiper deployment -> Multi-stage wiping"
        }
      ]
    },
    {
      "url": "https://www.zscaler.com/blogs/security-research/apt28-leverages-cve-2026-21509-operation-neusploit/",
      "created": "02/02/2026",
      "tags": [
        "APT28",
        "Operation Neusploit",
        "CVE-2026-21509",
        "MiniDoor",
        "PixyNetLoader",
        "Covenant",
        "Russia",
        "Ukraine",
        "Slovakia",
        "Romania"
      ],
      "iocs": [
        {
          "id": "ioc-1",
          "type": "file",
          "malign": true,
          "hash": "md5",
          "value": "f3b869a8d5ad243e35963ba6d7f89855",
          "comment": "Dropper DLL (Variant 1) filename: 2_2.d"
        },
        {
          "id": "ioc-2",
          "type": "file",
          "malign": true,
          "hash": "sha1",
          "value": "c1b272067491258ea4a2b1d2789d82d157aaf90a",
          "comment": "Dropper DLL (Variant 1) filename: 2_2.d"
        },
        {
          "id": "ioc-3",
          "type": "file",
          "malign": true,
          "hash": "sha256",
          "value": "a944a09783023a2c6c62d3601cbd5392a03d808a6a51728e07a3270861c2a8ee",
          "comment": "Dropper DLL (Variant 1) filename: 2_2.d"
        },
        {
          "id": "ioc-4",
          "type": "file",
          "malign": true,
          "hash": "md5",
          "value": "859c4b85ed85e6cc4eadb1a037a61e16",
          "comment": "PixyNetLoader dropper DLL (Variant 2) filename: table.d"
        },
        {
          "id": "ioc-5",
          "type": "file",
          "malign": true,
          "hash": "md5",
          "value": "e4a5c4b205e1b80dc20d9a2fb4126d06",
          "comment": "Shellcode loader filename: EhStoreShell.dll"
        },
        {
          "id": "ioc-6",
          "type": "file",
          "malign": true,
          "hash": "md5",
          "value": "154ff6774294e0e6a46581c8452a77de",
          "comment": "PNG container with embedded shellcode filename: SplashScreen.png"
        },
        {
          "id": "ioc-7",
          "type": "domain",
          "malign": true,
          "value": "freefoodaid.com",
          "comment": "Malicious domain used to host payloads"
        },
        {
          "id": "ioc-8",
          "type": "domain",
          "malign": true,
          "value": "wellnesscaremed.com",
          "comment": "Malicious domain used in campaign infrastructure"
        },
        {
          "id": "ioc-9",
          "type": "url",
          "malign": true,
          "value": "https://freefoodaid.com/documents/2_2.d",
          "comment": "URL hosting MiniDoor dropper DLL"
        },
        {
          "id": "ioc-10",
          "type": "url",
          "malign": true,
          "value": "https://freefoodaid.com/tables/tables.d",
          "comment": "URL hosting PixyNetLoader"
        },
        {
          "id": "ioc-11",
          "type": "process",
          "malign": false,
          "value": "explorer.exe",
          "comment": "EhStoreShell.dll executes malicious logic only when loaded by explorer.exe"
        },
        {
          "id": "ioc-12",
          "type": "cmdline",
          "malign": true,
          "value": "schtasks.exe /Create /tn \"OneDriveHealth\" /XML \"%temp%\\Diagnostics\\office.xml\"",
          "comment": "Scheduled task creation command used for staged execution"
        }
      ],
      "ttps": [
        {
          "id": "T1566.001",
          "comment": "Phishing: Spearphishing Attachment - weaponized RTFs delivered as attachments."
        },
        {
          "id": "T1203",
          "comment": "Exploitation for Client Execution - CVE-2026-21509 exploited by crafted RTF files."
        },
        {
          "id": "T1053.005",
          "comment": "Scheduled Task/Job: Scheduled Task - OneDriveHealth task created via schtasks.exe."
        },
        {
          "id": "T1546.015",
          "comment": "Event Triggered Execution: Component Object Model Hijacking - persistence via CLSID InProcServer32 hijack."
        },
        {
          "id": "T1137.006",
          "comment": "Office Application Startup: Add-ins - malicious Outlook VBA project loaded at startup."
        },
        {
          "id": "T1027.003",
          "comment": "Obfuscated Files or Information: Steganography - shellcode hidden in SplashScreen.png LSB data."
        },
        {
          "id": "T1497.003",
          "comment": "Virtualization/Sandbox Evasion: Time Based Checks - Sleep timing guardrail check in loader."
        },
        {
          "id": "T1114",
          "comment": "Email Collection - MiniDoor Outlook VBA forwards mailbox messages to attacker-controlled addresses."
        },
        {
          "id": "T1071.001",
          "comment": "Application Layer Protocol: Web Protocols - Covenant Grunt uses HTTPS for C2 traffic."
        },
        {
          "id": "T1102.002",
          "comment": "Web Service: Bidirectional Communication - Filen API abused as C2Bridge channel."
        },
        {
          "id": "T1480.002",
          "comment": "Execution Guardrails: Mutual Exclusion - campaign components create mutexes to avoid duplicate execution."
        }
      ],
      "queries": [
        {
          "id": "query-1",
          "type": "kql",
          "value": "DeviceProcessEvents | where FileName =~ 'schtasks.exe' | where ProcessCommandLine has '/Create' and ProcessCommandLine has 'OneDriveHealth' and ProcessCommandLine has 'office.xml'",
          "iocs": [
            "ioc-12"
          ],
          "comment": "Detect scheduled task creation observed in PixyNetLoader chain"
        },
        {
          "id": "query-2",
          "type": "kql",
          "value": "DeviceNetworkEvents | where RemoteUrl has_any ('freefoodaid.com', 'wellnesscaremed.com')",
          "iocs": [
            "ioc-7",
            "ioc-8",
            "ioc-9",
            "ioc-10"
          ],
          "comment": "Detect network activity to campaign infrastructure"
        },
        {
          "id": "query-3",
          "type": "yara",
          "value": "rule APT28_Operation_Neusploit_PixyNetLoader { strings: $m1 = \"asagdugughi41\" ascii $m2 = \"OneDriveHealth\" ascii $m3 = \"EhStoreShell.dll\" ascii condition: 2 of ($m*) }",
          "iocs": [
            "ioc-4",
            "ioc-5",
            "ioc-12"
          ],
          "comment": "Static detections based on strings explicitly documented in the report"
        }
      ],
      "xrefs": [
        {
          "id": "xref-1",
          "type": "detects",
          "parents": [
            "query-1",
            "query-2",
            "query-3"
          ],
          "children": [
            "ioc-7",
            "ioc-9",
            "ioc-10"
          ],
          "comment": "Operation Neusploit chain: exploit RTF -> dropper DLL -> persistence via COM hijack and scheduled task -> shellcode loader -> Covenant Grunt C2"
        }
      ]
    }
  ]
}